# conf/server/supervisor-app.conf
[program:app-uwsgi]
command = /usr/local/bin/uwsgi --ini /code/conf/server/uwsgi.ini

[program:nginx-app]
command = /usr/sbin/nginx
